# warehouse.service
# 仓储管理系统 systemd 服务文件
# 复制到 /etc/systemd/system/warehouse.service

[Unit]
Description=仓储管理系统 (Warehouse Management System)
Documentation=https://github.com/your-repo/warehouse-management
After=network.target mysql.service redis-server.service
Wants=mysql.service redis-server.service
Requires=network.target

[Service]
Type=notify
User=warehouse
Group=warehouse
WorkingDirectory=/opt/warehouse
Environment=FLASK_ENV=production
Environment=PYTHONPATH=/opt/warehouse
ExecStart=/opt/warehouse/venv/bin/gunicorn -c gunicorn_production.py app:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30
PrivateTmp=true
Restart=always
RestartSec=10

# 安全配置
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/warehouse /var/log/warehouse /var/backups/warehouse /tmp
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# 资源限制
LimitNOFILE=65536
LimitNPROC=32768

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=warehouse

[Install]
WantedBy=multi-user.target
